# wifi-analyzer.service
[Unit]
Description=WiFi Speed Analyzer

[Service]
Type=simple
ExecStart=/usr/bin/bash /home/<user>/wifi-speed-analyzer/analyzer.sh
